# debian/bakapy.upstart
description "Bakapy scheduler"

start on runlevel [2345]
stop on runlevel [016]

script
    BAKAPY_CONFIG=/etc/bakapy/bakapy.conf
    BAKAPY_LOGLEVEL=info
    BAKAPY_OPTS=

    if [ -f /etc/default/bakapy ]; then
        . /etc/default/bakapy
    fi

    exec /usr/bin/bakapy-scheduler --config=${BAKAPY_CONFIG} --loglevel=${BAKAPY_LOGLEVEL} $BAKAPY_OPTS
end script
